# supervisor_get_rich.conf
; ------------------------------------------------------------------------------
; ------------------------------------------------------------------------------
; ------------------------------------------------------------------------------
; SUPERVISOR
; ------------------------------------------------------------------------------
; ------------------------------------------------------------------------------
; ------------------------------------------------------------------------------

; STARTS THE BOT ON STARTUP
; RESTARTS HEN IT CRASHES
; http://supervisord.org/introduction.html

; ------------------------------------------------------------------------------
; INSTALLATION STEPS ON AMAZON AWS LINUX
; ------------------------------------------------------------------------------

; sudo amazon-linux-extras install epel
; sudo yum install supervisor

; ADD FOLLOWING SECTION TO CONFIG (etc/supervisord.conf)

[program:get-rich]
directory=/home/ec2-user/git/cryptocagibi
command=/usr/bin/python3.7 /home/ec2-user/git/cryptocagibi/get-rich-quick-scheme.py
autostart=true
autorestart=true
startretries=9
stderr_logfile=/home/ec2-user/git/cryptocagibi/supervisor.log
stdout_logfile=/home/ec2-user/git/cryptocagibi/supervisor.log
user=ec2-user

; ------------------------------------------------------------------------------
; FIX ERRORS AND WARNINGS
; ------------------------------------------------------------------------------

; Get rid of "unix:///tmp/supervisor.sock no such file" warning:
; --> sudo touch /var/run/supervisor.sock
; --> sudo chmod 777 /var/run/supervisor.sock

; ------------------------------------------------------------------------------
; USEFULL ACTIONS FOR DEBUG AND SETUP
; ------------------------------------------------------------------------------

; EDIT CONFIG FILE ---------------------> sudo vim /etc/supervisord.conf

; SUPERVISORD (/usr/bin/)
; STOP ------> sudo service supervisord stop
; START -----> sudo service supervisord start

; SUPERVISORCTL (command line client) (/usr/bin/)
; SHOW RUNNNG PROCESES --> sudo supervisorctl status

; SHOW RUNNING SUPERVISOR PROCESS ------> ps aux | grep superv

; LOGS:
; SHOW SUPERVISOR LOG --------------> sudo cat /var/log/supervisor/supervisord.log
; SHOW SUPERVISOR LOG II -----------> cat /tmp/supervisord.log
; SHOW PYTHON PRINT LOG-------------> cat ~/git/cryptocagibi/supervisor.out

